# overrides/etc/profile.d/99-penv.sh
export PATH="$PATH:/usr/local/games:/usr/games"